# autostart/dunst.desktop
[Desktop Entry]
Type=Application
Name=Dunst Notification provider
Exec=dunst
