# autostart/flameshot.desktop
[Desktop Entry]
Name=Flameshot
Comment=
Exec=flatpak run org.flameshot.Flameshot
Type=Application
