# autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=Network
Comment=Manage your network connections
Exec=nm-applet
