# autostart/picom.desktop
[Desktop Entry]
Name=Picom
Comment=
Exec=/usr/bin/picom
Type=Application
