# autostart/xdg-user-dirs.desktop
[Desktop Entry]
Type=Application
Name=User folders update
Exec=xdg-user-dirs-update
